# desktop/image/app/conf.d/xsession.conf
[program:xsession]
#command=/etc/X11/Xsession
command=jwm -f /app/system.jwmrc
user=developer
autorestart=true
stdout_logfile=/app/logs/xsession.log
redirect_stderr=true