# desktop/image/app/conf.d/xvnc.conf
[program:xvnc]
command=/app/bin/xvnc.sh
autorestart=true
stdout_logfile=/app/logs/xvnc.log
redirect_stderr=true